# prebuilt/common/etc/cron.d/root
* */2 * * * /system/bin/trim_partitions
